# owen-kbd-layout/etc/systemd/system/keymapping.service
[Unit]
Description=Runs the key mapper

[Service]
Type=simple
ExecStart=/usr/bin/map-all-keyboards

[Install]
WantedBy=multi-user.target
